# overlay/rootfs_overlay/etc/profile.d/monitorVariables.sh
#!/bin/bash
SETTINGS='/data/settings.txt'
NETWORK_CONF='/etc/systemd/network/wired.network'
MY_LIB='/usr/lib/monitor_lib.sh'
DATA='/data'