# overlay/rootfs_overlay/etc/profile.d/monitorVariables.sh
#!/bin/bash
SETTINGS='/data/settings.txt'
