# overlay/board/beaglebone/systemd/zabbixDumpTimer.service
[Unit]
Description=dump script for zabbix DB

[Service]
ExecStart=/usr/bin/zabbixDump.sh