# overlay/systemd/rodos.service
[Unit]
Description=RODOS temp daemon

[Service]
ExecStart=/usr/bin/rodos.sh